# device-support/Seeed-LTE_Cat_1_Pi_HAT/debian/postinst
#!/bin/bash

systemctl enable seeed_lte_hat
systemctl daemon-reload
systemctl start seeed_lte_hat
